# scripts/rc.local
#!/bin/bash

mkdir -p /metadata
mount LABEL=METADATA /metadata

/metadata/init.sh
/metadata/setnic.sh
/metadata/setssh.sh